# board/beamer/rootfs-overlay/etc/init.d/S49provision
#!/bin/sh

### BEGIN INIT INFO
# Provides:          beamer-provision
# Required-Start:    $local_fs $network
# Required-Stop:     $local_fs $network
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# Short-Description: Bring up provisioning SoftAP when not provisioned
### END INIT INFO

APP_DIR=/opt/beamer
APP_BIN="/usr/bin/python3 $APP_DIR/provision_app.py"
APP_PORT=80
DNSMASQ_PID=/var/run/dnsmasq-provision.pid
HOSTAPD_CONF_RUNTIME=/run/hostapd-provision.conf
NETPLAN_CONFIG="/boot/network-config"
LOG_TAG="S49provision"

# Common Beamer init helpers
[ -r /etc/init.d/common-func ] && . /etc/init.d/common-func

start() {
    # Run provisioning logic in the background so we don't block the rest of boot.
    # The behaviour is unchanged (wait for IPv4, then start AP/portal if needed),
    # but other services like sshd and the Beamer app can start in parallel.
    (
        # If wlan0 already has an IPv4 address, assume Wi‑Fi is up and skip AP.
        # This covers both cases where /boot/network-config was provided manually
        # and where it was generated by the portal.
        log_info "Checking existing Wi‑Fi connection on wlan0 (background)..."
        if [ -f "$NETPLAN_CONFIG" ]; then
            WAIT_SECS=20
            while [ "$WAIT_SECS" -gt 0 ]; do
                if ip -4 -o addr show dev wlan0 2>/dev/null | grep -q inet; then
                    log_info "wlan0 already has IPv4; skipping provisioning AP"
                    exit 0
                fi
                log_debug "wlan0 does not have IPv4; waiting for $WAIT_SECS seconds..."
                sleep 2
                WAIT_SECS=$((WAIT_SECS - 2))
            done
        else
            log_info "No $NETPLAN_CONFIG; skipping IPv4 wait and starting provisioning AP immediately"
        fi

        log_warn "Starting provisioning AP on wlan0 (no IPv4 after timeout)..."

        # Ensure we are not fighting with any client Wi‑Fi processes (e.g. wpa_supplicant/dhclient)
        # that might try to reconfigure wlan0 back into station mode.
        if pgrep -x wpa_supplicant >/dev/null 2>&1; then
            log_info "Stopping wpa_supplicant for provisioning AP..."
            killall wpa_supplicant 2>/dev/null || true
        fi
        if pgrep -x dhclient >/dev/null 2>&1; then
            log_info "Stopping dhclient for provisioning AP..."
            killall dhclient 2>/dev/null || true
        fi

        ip link set wlan0 up 2>/dev/null
        ip addr flush dev wlan0
        ip addr add 192.168.4.1/24 dev wlan0

        # Prepare runtime hostapd config using the current hostname as SSID
        mkdir -p "$(dirname "$HOSTAPD_CONF_RUNTIME")"
        SSID_NAME=$(hostname 2>/dev/null || echo "beamer")
        cat >"$HOSTAPD_CONF_RUNTIME" <<EOF
interface=wlan0
driver=nl80211
ssid=$SSID_NAME
hw_mode=g
channel=6
wmm_enabled=1
ignore_broadcast_ssid=0
EOF

        # Start dnsmasq and hostapd
        if dnsmasq --conf-file=/etc/dnsmasq.conf --pid-file=$DNSMASQ_PID; then
            log_info "dnsmasq started with /etc/dnsmasq.conf (pid file $DNSMASQ_PID)"
        else
            log_err "Failed to start dnsmasq with /etc/dnsmasq.conf"
        fi

        if hostapd -B "$HOSTAPD_CONF_RUNTIME"; then
            log_info "hostapd started with runtime config $HOSTAPD_CONF_RUNTIME (SSID: $SSID_NAME)"
        else
            log_err "Failed to start hostapd with runtime config $HOSTAPD_CONF_RUNTIME"
        fi

        # Start Flask app on port 80 if not already running
        if ! pgrep -f "$APP_BIN" >/dev/null 2>&1; then
            log_info "Starting portal app on port $APP_PORT..."
            APP_PORT=$APP_PORT $APP_BIN &
        fi
    ) &

    return 0
}

stop() {
    log_info "Stopping provisioning AP on wlan0..."
    killall hostapd 2>/dev/null || true
    if [ -f "$DNSMASQ_PID" ]; then
        kill "$(cat $DNSMASQ_PID)" 2>/dev/null || true
        rm -f "$DNSMASQ_PID"
    else
        killall dnsmasq 2>/dev/null || true
    fi
}

case "$1" in
  start)
    start
    ;;
  stop)
    stop
    ;;
  restart)
    stop
    start
    ;;
  *)
    echo "Usage: $0 {start|stop|restart}"
    exit 1
    ;;
esac

exit 0


